# Chapter01/goproject.conf
[supervisord]
logfile = /tmp/supervisord.log

[program:myserver]
command=/home/floyed/golab/bin/myserver
autostart=true
autorestart=true
redirect_stderr=true
